# src/modules/wifipi/filesystem/root/etc/systemd/system/autohotspot.timer
[Unit]
Description=Run the AutoHotspot script once per minute

[Timer]
OnUnitInactiveSec=1min

[Install]
WantedBy=timers.target